# rootfs-main_20121226⁄/usr/lib/opkg/info/busybox.list
/etc/busybox.links
/bin/sh
/etc/udhcpc.d/50default
/usr/share/udhcpc/default.script
/bin/busybox
/etc/default/hwclock
/etc/init.d/hwclock.sh
