# rootfs-main_20121226⁄/usr/lib/opkg/info/initscripts.list
/etc/default/devpts
/etc/init.d/sysfs.sh
/etc/init.d/umountfs
/etc/rcS.d/S30ramdisk
/etc/rc0.d/S20sendsigs
/etc/rc0.d/S40umountfs
/etc/rcS.d/S38devpts.sh
/etc/rcS.d/S45mountnfs.sh
/etc/default/volatiles/00_core
/etc/init.d/ramdisk
/etc/rcS.d/S10checkroot
/etc/init.d/populate-volatile.sh
/etc/rc4.d/S99rmnologin
/etc/init.d/bootmisc.sh
/etc/init.d/checkroot
/etc/rcS.d/S03sysfs
/etc/rc2.d/S99rmnologin
/etc/rc6.d/S31umountnfs.sh
/etc/init.d/devpts.sh
/etc/rc0.d/S90halt
/etc/init.d/urandom
/etc/rc0.d/S25save-rtc.sh
/etc/init.d/halt
/etc/init.d/banner
/etc/init.d/hostname.sh
/etc/init.d/sendsigs
/etc/rcS.d/S06alignment
/etc/rcS.d/S37populate-volatile.sh
/etc/init.d/save-rtc.sh
/etc/rc6.d/S25save-rtc.sh
/etc/rc5.d/S99rmnologin
/etc/rc6.d/S20sendsigs
/etc/rc6.d/S40umountfs
/etc/rcS.d/S39hostname.sh
/etc/rc6.d/S90reboot
/etc/init.d/reboot
/etc/rc3.d/S99rmnologin
/etc/rcS.d/S02banner
/etc/rcS.d/S35mountall.sh
/etc/init.d/alignment.sh
/etc/device_table
/etc/init.d/umountnfs.sh
/etc/init.d/functions
/etc/init.d/rmnologin
/etc/rcS.d/S55bootmisc.sh
/etc/init.d/mountnfs.sh
/etc/rcS.d/S99finish.sh
/etc/init.d/mountall.sh
/etc/rc0.d/S31umountnfs.sh
/etc/init.d/single
/etc/init.d/finish.sh
