# rootfs-main_20121226⁄/usr/lib/opkg/info/netbase.list
/usr/sbin/update-inetd
/etc/protocols
/etc/init.d/networking
/etc/network/interfaces
/etc/services
/etc/network/options
/etc/hosts
/etc/rpc
