# rootfs-main_20121226⁄/usr/lib/opkg/info/netbase.postinst
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT="-s"
fi
update-rc.d $OPT networking start 41 S . stop 41 0 6 1 .
