# rootfs-main_20121226⁄/etc/profile.d/qt.sh
#!/bin/sh

export QWS_MOUSE_PROTO=tslib:/dev/ttyO3
