# rootfs-main_20121226⁄/etc/profile.d/tslib.sh
#!/bin/sh

export TSLIB_CALIBFILE=/etc/pointercal
export TSLIB_CONFFILE=/etc/ts.conf
export TSLIB_PLUGINDIR=/usr/lib/ts/

export TSLIB_TSDEVICE=/dev/ttyO3
